# systemd/.config/systemd/user/get_mail.service
[Unit]
Description=Run getmail

[Service]
ExecStart=/usr/bin/getmail --quiet

[Install]
WantedBy=default.target
